# tmate/session@.service
[Unit]
Description=Starts tmate session
PartOf=workers.target

[Service]
Type=simple
ExecStartPre=/usr/bin/install -m 777 -d /tmp/devices
ExecStartPre=/bin/su villager-%i -c "install -m 700 -d /tmp/devices/%i"
ExecStartPre=-/usr/bin/pkill -9 -u villager-%i
ExecStart=/bin/su villager-%i -c "/usr/bin/tmate -S /tmp/devices/%i/%i.sock new-session"
EnvironmentFile=/root/%i
Restart=on-failure

[Install]
WantedBy=multi-user.target
