# tmate/session@.service
[Unit]
Description=Starts tmate session
PartOf=workers.target

[Service]
EnvironmentFile=/root/data
User=villager
Type=forking
ExecStart=/root/session.sh %i
Restart=on-failure
ExecReload=/root/session_restart.sh %i

[Install]
WantedBy=multi-user.target
