# tmate/session@.service
[Unit]
Description=Starts tmate session
PartOf=workers.target

[Service]
User=villager
Type=simple
ExecStartPre=/bin/mkdir -p /tmp/devices/%i
ExecStart=/usr/bin/tmate -F -S /tmp/devices/%i/%i.sock new-session
EnvironmentFile=/root/device%i
Restart=on-failure

[Install]
WantedBy=multi-user.target
